# OpenStack/misc/cloudresources.conf
description "Start and stop the cloud resources service"

env BIN_DIR=/usr/local/bin
env PORT=8080

start on startup

exec $BIN_DIR/cloudresources $PORT
